# renterd_1.0.2_arm64/DEBIAN/prerm
#!/bin/sh
systemctl stop renterd.service
systemctl disable renterd.service
